# osquery.conf
{
  "options": {
    "host_identifier": "hostname",
    "schedule_splay_percent": 10
  },
  "scheduledQueries": [

    {
      "name": "LaunchD modified",
      "query": "SELECT h.path, h.sha1 FROM launchd as l join hash as h where h.path=l.program || l.program_arguments and run_at_load = 'true' and keep_alive = 'true';", 
      "interval": 60
    },

    {
      "name": "Startup Items modified",
      "query": "SELECT name FROM startup_items",
      "interval": 60
    },
    {
      "name": "Crontab was modified",
      "query": "SELECT * FROM crontab",
      "interval": 60
    },
    {
      "name": "Etc_hosts was modified",
      "query": "SELECT * FROM etc_hosts",
      "interval": 60
    },
    {
      "name": "Groups were modified",
      "query": "SELECT * FROM groups",
      "interval": 60
    },
    {
      "name": "Interfaces were modified",
      "query": "SELECT * FROM interface_addresses",
      "interval": 60
    },
    {
      "name": "Kernel information was modified",
      "query": "SELECT * FROM kernel_info",
      "interval": 60
    },
    {
      "name": "Passwd was changed",
      "query": "SELECT * FROM passwd",
      "interval": 60
    },
    {
      "name": "SUID binaries changed",
      "query": "SELECT * FROM suid_bin",
      "interval": 60
    },
    {
      "name": "Users Changed",
      "query": "SELECT * FROM users",
      "interval": 60
    },
    {
      "name": "App Layer Firewall Changed",
      "query": "SELECT * FROM alf",
      "interval": 60
    },
    {
      "name": "App Layer Firewall Exception Change",
      "query": "SELECT * FROM alf_exceptions",
      "interval": 60
    },
    {
      "name": "App Layer Firewall Explicit Auths Change",
      "query": "SELECT * FROM alf_explicit_auths",
      "interval": 60
    },
    {
      "name": "Application change",
      "query": "SELECT * FROM apps",
      "interval": 60
    },
    {
      "name": "Certificate Authority Change",
      "query": "SELECT * FROM ca_certs",
      "interval": 60
    },
    {
      "name": "Something added to Quarantine",
      "query": "SELECT * FROM quarantine",
      "interval": 60
    },
    {
      "name": "Xprotect report happened",
      "query": "SELECT * FROM xprotect_reports",
      "interval": 60
    }
  ],
  
  "additional_monitoring" : {
    "file_paths": {
      "system_files": [
    "/etc"
      ],
      "core_binaries": [
        "/bin/"
      ],
      "other_binaries":[
        "/sbin/",
    "/usr/sbin/%"
      ]
    }
  }
}
